# services/vane/vane-ai.service
[Unit]
Description=Vane AI Search Engine — SearxNG + Groq Research Platform
Documentation=file:///home/aseps/MCP/README.md
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=aseps
RemainAfterExit=yes

ExecStart=/usr/bin/docker start vane
ExecStop=/usr/bin/docker stop vane

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vane-ai

[Install]
WantedBy=multi-user.target
